# centos7rails/etc/supervisord.d/mariadb.ini
[program:mariadb]
command=/etc/init.d/mysql start
autostart=true
autorestart=true
