# centos7rails/etc/supervisord.d/redis.ini
[program:redis]
command=/usr/bin/redis-server /etc/redis.conf
autstart=true
autrestart=true

