# postfix/2.10/postfix.ini
[program:postfix]
command = /usr/sbin/postfix -c /etc/postfix start
process_name = master
directory = /etc/postfix
startsecs = 0
autorestart = false
